# packages/platforms/surfacert/retroarch-config-surfacert/retroarch-config-surfacert.install
post_install() {
  # enable automatic login to xorg/retroarch 
  systemctl enable xlogin@root
}

