# rr-packages/system/rr-overlay-surfacert/overlay/etc/systemd/system/surfacert-mlan0.service
[Unit]
Description=surfacert mlan0 connection (wifi)

[Service]
ExecStart=/usr/bin/srtwconnect

[Install]
WantedBy=multi-user.target

