# overlays/common/etc/systemd/system/retroroot.service
[Unit]
Description=RetroRoot data partition resizer

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/rr-resize

